# docker-pathfinder/supervisord-pathfinder.conf
[supervisord]
nodaemon=true

[program:php]
command=/usr/sbin/php-fpm7.0 -F
priority=3
autostart=true
autorestart=true

[program:nginx]
command=/usr/sbin/nginx
stdout_events_enabled=true
stderr_events_enabled=true
priority=4
autostart=true
autorestart=true
